# board/raspberrypi4/rootfs-overlay/etc/profile.d/locale.sh
#!/bin/sh

export LANG=uk_UA.UTF-8
export LC_ALL=uk_UA.UTF-8
